# web/scripts/tiktok-studio/com.flashflow.tiktok-healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flashflow.tiktok-healthcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/brandonglomski/tts-engine/web/scripts/tiktok-studio/scheduled-healthcheck.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/brandonglomski/tts-engine/web</string>

    <!-- Run every 30 minutes (1800 seconds) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/brandonglomski/tts-engine/web/data/sessions/logs/healthcheck.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/brandonglomski/tts-engine/web/data/sessions/logs/healthcheck-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/brandonglomski</string>
    </dict>

    <!-- Don't restart on failure — next interval will retry -->
    <key>KeepAlive</key>
    <false/>

    <!-- Allow macOS to delay up to 5 minutes if system is busy -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
